# ai-frontend-communication/ops/systemd/dca-api.service
[Unit]
Description=Developer Communication Agent API
After=network-online.target
Wants=network-online.target
ConditionFileIsExecutable=/opt/dca/current/.venv/bin/uvicorn

[Service]
Type=simple
User=dca
Group=dca
WorkingDirectory=/opt/dca/current
EnvironmentFile=/etc/dca/dca.env
Environment=HOME=/var/lib/dca
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/dca/current/.venv/bin/uvicorn dca.app:app --host 172.18.0.1 --port 8000 --proxy-headers --forwarded-allow-ips=172.18.0.0/16
Restart=on-failure
RestartSec=3
TimeoutStopSec=30
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/lib/dca

[Install]
WantedBy=multi-user.target
